# .codex/skills-archive/debian-web-stack/assets/systemd-unit.service
[Unit]
Description=example web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=example
Group=example
WorkingDirectory=/opt/example
EnvironmentFile=/etc/example/example.env
ExecStart=/usr/bin/env bash -lc 'node server.js'
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

