# .codex/skills-archive/linux-timers-and-cron/assets/example.service
[Unit]
Description=Example scheduled job (service)

[Service]
Type=oneshot
WorkingDirectory=/opt/job
ExecStart=/usr/bin/env bash -lc '/opt/job/run.sh'
